# packaging/alpine/birdhub.initd
#!/sbin/openrc-run

name="birdhub"
description="VISP (Virtual ISP) Zero-Trust Overlay Network Agent"

command="/usr/bin/birdhub"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
directory="/etc/birdhub"

depend() {
    need net localmount
    use netbird
    after netbird firewall
}

start_pre() {
    # Ensure the configuration directory exists
    checkpath -d -m 0755 -o root:root /etc/birdhub
}
